# beerboard.service
[Unit]
Description=BeerBoard

[Service]
ExecStart=/opt/beerboard/beerboard
WorkingDirectory=/opt/beerboard


[Install]
WantedBy=multi-user.target
